# albert/.config/autostart/albert.desktop
[Desktop Entry]
Type=Application
Name=Albert
Exec=albert
Icon=albert
